# etc/pop-os/update-motd.d/10-help-text
#!/bin/sh

echo
echo " * Homepage: https://pop.system76.com"
echo " * Support:  https://support.system76.com"
echo
